# nginx/init/nginx-conf-watcher
### BEGIN INIT INFO
# Provides:          nginx-conf-watcher
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Monitor Nginx conf file changes
# Description:       Check changes into nginx.conf file and reload Nginx server.
### END INIT INFO

#!/bin/bash

USER=root

case "$1" in
        start)
            echo "Start Nginx config file monitoring"
            cd /usr/local/nginx/conf-watcher && sudo -u ${USER} forever start server.js
            ;;
        status)
            sudo -u ${USER} forever list | grep server.js | grep -v STOPPED
            ;;
        stop)
            echo "Stop Nginx config file monitoring"
            cd /usr/local/nginx/conf-watcher && sudo -u ${USER} forever stop server.js
            ;;
        *)
            echo "Usage: /etc/init.d/nginx-conf-watcher {start|stop|status}" > /dev/stderr
            exit 1
        ;;
esac